# server/supervisor/local.ini
[program:fabrydb-local]
directory = /home/devuser/.virtualenvs/fabrydb-env/fabrydb
command = /home/devuser/.virtualenvs/fabrydb-env/bin/uwsgi --ini /home/devuser/.virtualenvs/fabrydb-env/fabrydb/server/uwsgi/local.ini
autostart = true
autorestart = true
