# server/supervisor/production.ini
[program:fabrydb-production]
directory = /home/devuser/webapps/fabrydb-env/fabrydb
command = /home/devuser/webapps/fabrydb-env/bin/uwsgi --ini /home/devuser/webapps/fabrydb-env/fabrydb/server/uwsgi/production.ini
autostart = true
autorestart = true
